# m-tunnel.service
[Unit]
Description=M-Tunnel SSH Tunnel Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/m-tunnel-rust
Restart=always
RestartSec=5
User=m-tunnel
Group=m-tunnel
WorkingDirectory=/etc/m-tunnel
Environment=RUST_LOG=info

# Enhanced Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
PrivateTmp=true
PrivateDevices=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
ReadWritePaths=/var/log/m-tunnel /tmp
RemainAfterExit=false

[Install]
WantedBy=multi-user.target
