# m-tunnel-rust_1.0.0-1/lib/systemd/system/m-tunnel.service
[Unit]
Description=A Rust-based SSH tunnel utility for forward and reverse port forwarding
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/m-tunnel-rust
Restart=always
RestartSec=5
User=root
Group=root
WorkingDirectory=/etc/m-tunnel
Environment=RUST_LOG=info

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/m-tunnel
PrivateTmp=true

[Install]
WantedBy=multi-user.target
