# m-tunnel-rust_1.0.0-1/DEBIAN/postrm
#!/bin/bash
set -e
systemctl stop m-tunnel-rust.service || true
systemctl disable m-tunnel-rust.service || true
systemctl daemon-reload
