# m-tunnel-rust_1.0.0-1/DEBIAN/postinst
#!/bin/bash
set -e
systemctl daemon-reexec || true
systemctl daemon-reload
systemctl enable m-tunnel-rust.service
systemctl start m-tunnel-rust.service
