# m-tunnel-rust_1.0.0-1/DEBIAN/postinst
#!/bin/bash
set -e

# Create dedicated user for security
if ! id "m-tunnel" &>/dev/null; then
    useradd -r -s /bin/false -d /etc/m-tunnel -c "M-Tunnel Service User" m-tunnel
    echo "Created m-tunnel user"
fi

# Create and secure directories
mkdir -p /var/log/m-tunnel
mkdir -p /etc/m-tunnel
touch /etc/m-tunnel/known_hosts

# Set proper ownership and permissions
chown -R m-tunnel:m-tunnel /etc/m-tunnel
chown -R m-tunnel:m-tunnel /var/log/m-tunnel

# Secure permissions
chmod 750 /etc/m-tunnel
chmod 750 /var/log/m-tunnel
chmod 644 /etc/m-tunnel/m-tunnel.conf 2>/dev/null || true
chmod 600 /etc/m-tunnel/.env 2>/dev/null || true
chmod 600 /etc/m-tunnel/m-tunnel.key 2>/dev/null || true
chmod 644 /etc/m-tunnel/known_hosts

# Allow m-tunnel user to create SSH control sockets
mkdir -p /tmp/ssh-m-tunnel
chown m-tunnel:m-tunnel /tmp/ssh-m-tunnel
chmod 700 /tmp/ssh-m-tunnel

# Reload systemd and enable service
systemctl daemon-reload
systemctl enable m-tunnel.service

echo "=== M-Tunnel Installation Complete ==="
echo "Configuration files are located in /etc/m-tunnel/"
echo "1. Edit /etc/m-tunnel/.env with your SSH settings"
echo "2. Edit /etc/m-tunnel/m-tunnel.conf with tunnel definitions"  
echo "3. Add SSH host key: ssh-keyscan -H your-ssh-host >> /etc/m-tunnel/known_hosts"
echo "4. Start the service: systemctl start m-tunnel.service"
echo ""
echo "Security: Service runs as dedicated 'm-tunnel' user"
echo "Logs: journalctl -u m-tunnel -f"
